# scan_wifi.service
[Unit]
Description=Scan Wifi
After=multi-user.target

[Service]
WorkingDirectory=/home/julien/git/Scrin-Scanner
ExecStart=/usr/bin/sudo scan_wifi.py

[Install]
WantedBy=multi-user.target