# scan_wifi.service
[Unit]
Description=Scan Wifi
Wants=network-online.target
After=network.target network-online.target

[Service]
ExecStart=/usr/bin/sudo /usr/bin/python -u /home/julien/git/Scrin-Scanner/scan_wifi.py

[Install]
WantedBy=multi-user.target